# frontend.ini
[program:frontend]
user=root
environment=HOME="/root",USER="root"
directory=/root/lctc-ceph-performance-workbench/frontend
command=/usr/bin/npm start
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/root/frontend.out.log
stderr_logfile=/root/frontend.err.log
