# job_conductor.ini
[program:job_conductor]
user=root
environment=HOME="/root",USER="root"
directory=/root/lctc-ceph-performance-workbench/backend/ceph_perf_api/
command=/root/lctc-ceph-performance-workbench/backend/ceph_perf_api/run_job_conductor.py
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/root/job_conductor.out.log
stderr_logfile=/root/job_conductor.err.log
